# deploy/systemd/webfetch-maintenance.timer
[Unit]
Description=Run WebFetch maintenance hourly

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target

